# systemd/parallama.service
[Unit]
Description=Parallama API Gateway Service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=parallama
Group=parallama
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/parallama/environment
ExecStart=/usr/bin/parallama-cli serve --config /etc/parallama/config.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
